# infra/systemd/powerhouse-log-export@.service
[Unit]
Description=Power-House Log Export (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/powerhouse/powerhouse-common.env
EnvironmentFile=/etc/powerhouse/powerhouse-%i.env
ExecStart=/usr/local/lib/powerhouse/journal_export.sh
