# infra/systemd/powerhouse-healthcheck@.service
[Unit]
Description=Power-House Healthcheck (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/powerhouse/powerhouse-common.env
EnvironmentFile=/etc/powerhouse/powerhouse-%i.env
ExecStart=/usr/local/lib/powerhouse/healthcheck.py
