# infra/systemd/powerhouse-reliability-campaign.service
[Unit]
Description=Power House external 72-hour reliability campaign
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/systemd-inhibit --what=sleep:shutdown:idle --mode=block --why=Power-House-reliability-campaign /usr/bin/python3 %h/.local/lib/powerhouse/reliability_campaign.py run --config %h/.config/powerhouse/reliability-campaign.json
Restart=on-failure
RestartSec=15
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/state/powerhouse/reliability
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=default.target
