# infra/systemd/powerhouse-backup@.service
[Unit]
Description=Power-House Backup (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/jrocnet/powerhouse-common.env
EnvironmentFile=/etc/jrocnet/powerhouse-%i.env
ExecStart=/usr/local/lib/powerhouse/backup.sh
