# infra/monitoring/powerhouse-observer-registry.timer
[Unit]
Description=Refresh Power House public observer registry health

[Timer]
OnBootSec=25s
OnUnitActiveSec=15s
AccuracySec=1s
Persistent=true

[Install]
WantedBy=timers.target
