# powerhouse-boot2.service
[Unit]
Description=JROC-NET Boot2 (julian)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/lib/jrocnet/boot2
ExecStart=/usr/local/bin/julian net start \
  --node-id boot2 \
  --log-dir /var/lib/jrocnet/boot2/logs \
  --listen /ip4/0.0.0.0/tcp/7002 \
  --bootstrap /ip4/<boot1-ip>/tcp/7001/p2p/<boot1-peer-id> \
  --broadcast-interval 5000 \
  --quorum 2 \
  --key /etc/jrocnet/boot2.ed25519 \
  --policy /etc/jrocnet/governance.json \
  --checkpoint-interval 60
Restart=always
RestartSec=3
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
